# roles/mesh_fastd/templates/fastd.service.j2
[Unit]
Description=Aggregation of all fastd instances

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
