# roles/archive/ffh.crowdsec/files/crowdsec-update.timer
[Unit]
Description=Daily run crowdsec updates

[Timer]
OnCalendar=daily
RandomizedDelaySec=14400

[Install]
WantedBy=timers.target
