# roles/ffh.gw-statuspage/templates/gw-statuspage-restart.service.j2
[Unit]
Description=GW Statuspage Restart service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart gw-statuspage.service
