# roles/ffh.bird2/files/bird-peer-info.service
[Unit]
Description=Extracts route data from birdc for Zabbix to pick up

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bird-peer-info-routes.sh
