# roles/ffh.zabbix-sql/templates/mysql-zbx-part.service
[Unit]
Description=Partition zabbix database service
[Service]
Type=oneshot
ExecStart=/usr/lib/zabbix/mysql_zbx_part.pl
