# roles/ffh.zabbix-sql/templates/zabbix-backup.timer
[Unit]
Description=Backup Zabbix mysql db
[Timer]
OnCalendar=daily
[Install]
WantedBy=timers.target
