# roles/nginx/templates/08-getssl.cron.j2
#!/bin/bash
{{ nginx_getssl_dir }}/getssl -w {{ nginx_getssl_dir }} -U -a -q
