# roles/nginx/templates/getssl.cron.j2
#!/bin/bash
/opt/getssl/getssl -U -a -q
