# roles/ffh.apinger/files/apinger.service
[Unit]
Description=Apinger

[Service]
ExecStart=/usr/local/sbin/apinger -c /etc/apinger.conf -f

[Install]
WantedBy=multi-user.target
