# roles/ffh.restic/templates/restic-prune.timer
[Unit]
Description=Prune data from the restic repository monthly
[Timer]
OnCalendar=*-*-1 4:00:00
[Install]
WantedBy=timers.target
