# roles/ffh.restic/templates/restic-prune.timer
[Unit]
Description=Prune data from the restic repository monthly
[Timer]
OnCalendar=monthly
[Install]
WantedBy=timers.target
