# roles/ffh.restic/templates/restic-backup.timer
[Unit]
Description=Backup with restic daily
[Timer]
OnCalendar=daily
[Install]
WantedBy=timers.target
