# roles/ffh.supernode/templates/rc.local.j2
#!/bin/sh -e
#
# do not change this, since it is created by ansible!
#

#
# RULES FOR TABLE freifunk
# Moved to systemd-networkd
#

exit 0
