# roles/archive/node_exporter/templates/bpfcountd.service.j2
[Unit]
Description=Aggregation of all bpfcountd instances

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

